# build/config/blogserver.service

[Unit]
//Description=BlogNx Server

[Service]
//ExecStart=/usr/local/bin/node /var/www/blogserver/src/server.js
//Restart=on-failure
//Environment=PATH=/usr/bin:/usr/local/bin
//Environment=API_URL=https://www.blognx.com/api
//Environment=PORT=3000
//WorkingDirectory=/var/www/blogserver

[Install]
//WantedBy=multi-user.target